# init/horebd.socket
[Unit]
Conflicts=horebd.service

[Socket]
ListenStream=9999
Accept=yes

[Install]
WantedBy=sockets.target
