# init/horebd@.service
[Unit]
Description=Horebd Per-Connection Daemon

[Service]
ExecStart=-/usr/bin/horebd
KillMode=process
